# .config/autostart/Picom - Compositor.desktop
[Desktop Entry]
Exec=picom --daemon --backend glx
Name=Picom - Compositor
Type=Application
Version=1.0
