# .config/autostart/Picom - Compositor.desktop
[Desktop Entry]
Exec=picom -b
Name=Picom - Compositor
Type=Application
Version=1.0
